# updater/crontab
*/10 * * * * cd /app && /usr/local/bin/python src/main.py >> /var/log/pldashboard-updater.log 2>&1
