# updater/crontab
*/10 * * * * cd /app && /usr/local/bin/uv run updater >> /var/log/pldashboard-updater.log 2>&1
